# discord-bot.service
[Unit]
Description=Discord to Telegram Bot
After=network.target

[Service]
Type=exec
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/discord-telegram-bot
Environment=PATH=/home/ubuntu/discord-telegram-bot/venv/bin:/usr/bin:/bin
ExecStart=/usr/bin/python3 app.py
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=discord-bot

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/ubuntu/discord-telegram-bot /var/log/discord-bot

# Environment variables
Environment="FLASK_ENV=production"
Environment="PYTHONPATH=/home/ubuntu/discord-telegram-bot"
Environment="DISCORD_CLIENT_ID=1411364572274888836"
Environment="DISCORD_CLIENT_SECRET=Bd8MwfHxqfu0CdFE231z1Z4AG6aGE52"
Environment="DISCORD_REDIRECT_URI=https://13.60.40.184/discord/callback"
Environment="SECRET_KEY=a1b2c3d4e5f6789012345678901234567890abcdef1234567890abcdef123456"
Environment="FLASK_DEBUG=False"

[Install]
WantedBy=multi-user.target
